# .config/shell/profile.d/015-steam.sh
## ek9/dotfiles - https://github.com/ek9/dotfiles
## 015-steam.sh
## fixes missing menu text in civilization 5 and cs:go on linux
#export MESA_GL_VERSION_OVERRIDE=3.2
#export MESA_GLSL_VERSION_OVERRIDE=150

