# .config/shell/profile.d/03-dircolors.sh
## ek9/shell-config - https://github.com/ek9/shell-config
## 03-dircolors.sh
## Setup dircolors

# dircolors
[[ -f "$HOME/.config/dir_colors" ]] && eval $(dircolors --sh ~/.config/dir_colors)

