# .config/shell/profile.d/30-wine.sh
## ek9/shell-config - https://github.com/ek9/shell-config
## 30-wine.sh
## wine tweaks to force 32bit mode

# wine force 32bit mode
export WINEARCH=win32
