# .config/bash/profile.d/003-stty-tweak.sh
# disable Ctrl-S / Ctrl-Q combination which locks terminal scroll / vim
stty -ixon

