# .config/bash/profile.d/103-bash-theme.sh
# base16-shell themes
BASE16_SHELL="$HOME/var/cache/base16-shell.git/base16-default.dark.sh"
[[ -s $BASE16_SHELL ]] && source $BASE16_SHELL
