# .config/bash/profile.d/45-oh-my-git.sh
if [[ -f $HOME/.config/bash/oh-my-git.git/prompt.sh ]]; then
    source $HOME/.config/bash/oh-my-git.git/prompt.sh
fi
