# .config/bash/profile.d/45-oh-my-git.sh
if [[ -f $HOME/var/cache/mr/oh-my-git.git/prompt.sh ]]; then
    source $HOME/var/cache/mr/oh-my-git.git/prompt.sh
fi
